# systemd/framework2.service
[Unit]
Description=framework2 Timers Service
PartOf=framework2.target
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/var/www/framework2/scripts/timers.php
Restart=on-failure
RestartSec=5

[Install]
WantedBy=framework2.target
